# cron.conf
[program:cron]
command=/bin/bash -c "exec /usr/sbin/cron -f"
